# packages/debian/DEBIAN/postinst
#!/bin/sh

set -e

systemctl daemon-reload
systemctl enable nx-proxy.service
systemctl start nx-proxy.service

echo "nx-proxy service installed and started."
